# shock-os-mate-desktop-10-all/etc/systemd/user/shock-slideshow-background-daemon.service
[Unit]
Description=Slideshow Background Daemon

[Service]
Type=simple
ExecStart=/usr/bin/shock-slideshow-background-daemon

[Install]
WantedBy=graphical-session.target
