# shock-os-core-9-all/etc/systemd/system/shock-screentime-daemon.service
[Unit]
Description=Screen Time Counter
After=default.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/shock-screentime-daemon

[Install]
WantedBy=graphical.target
