# shock-os-core-10-all/etc/systemd/user/shock-updates-scanner.service
[Unit]
Description=Checks for system updates and notifies the user if any are available (or installs them automatically if automatic updates are enabled)

[Service]
Type=simple
ExecStart=/usr/bin/shock-updates-scanner

[Install]
WantedBy=graphical-session.target
