# shock-os-mate-desktop-9-all/etc/systemd/user/shock-notifications-log-daemon.service
[Unit]
Description=Notification Log Daemon

[Service]
Type=simple
ExecStart=/usr/bin/shock-notifications-log-daemon
