# shockos-9.0-all/etc/systemd/system/shock-slideshow-background-daemon.service
[Unit]
Description=Slideshow Background Daemon
After=graphical.target

[Service]
ExecStart=/usr/bin/shock-slideshow-background-daemon

[Install]
WantedBy=default.target
