# shockos-9.0-all/etc/systemd/system/shock-startup-daemon.service
[Unit]
Description=Shock Startup Daemon
After=graphical.target

[Service]
ExecStart=/usr/bin/shock-startup-daemon

[Install]
WantedBy=default.target
