# shockos-8.0-all/etc/xdg/autostart/startup.desktop
[Desktop Entry]
Name=startup.desktop
Exec=shock-startup
Name[en_US]=startup.desktop
X-GNOME-Autostart-Phase=Desktop
Type=Application
Terminal=false
