# shockos-9.0-all/etc/systemd/user/shock-update-manager.service
[Unit]
Description=Update Manager

[Service]
Type=simple
ExecStart=/usr/bin/shock-update-manager --systemd
