# shockos-9.0-all/etc/systemd/user/shock-release-scanner.service
[Unit]
Description=Check for upgrades to newer versions of Shock OS

[Service]
Type=simple
ExecStart=/usr/bin/shock-release-scanner
